# systemd/zulip-aws-usage.timer
[Unit]
Description=Periodic AWS Usage notification

[Timer]
OnCalendar=*-*-* 18:00:00
Persistent=true

[Install]
WantedBy=timers.target
